# recyclarr-bin/recyclarr.timer
[Unit]
Description=Daily Recyclarr Sync

[Timer]
OnCalendar=daily
RandomizedDelaySec=360min
Persistent=true

[Install]
WantedBy=timers.target
